# dashboard/conf/supervisor.conf
[program:dashboard]
command=/var/src/virtualenv/estivador/bin/uwsgi --ini /var/src/projects/estivador/dashboard/conf/uwsgi.ini
directory=/var/src/projects/estivador/dashboard/
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/estivador/supervisor-dashboard.log
